# scripts/supervisor.conf
[program:ReactWebNative8Koa]
command = node index.js 1>log/node.log 2>log/node.log 0>/log/node.log &
directory = /home/www/bin/ReactWebNative8Koa
user = www